# supervisord.conf
; =============================================================================
; supervisord.conf — Sovereign Pipeline Process Manager
; =============================================================================
;
; Manages the Ralph ceremonies delivery machine as a supervised process.
;
; USAGE:
;   scripts/supervisord-start.sh    # start supervisord + managed processes
;   scripts/supervisord-stop.sh     # graceful stop
;   supervisorctl status            # check process state
;   supervisorctl tail -f ralph     # stream logs
;
; SECURITY: Change the web UI password before exposing this to any network.
; See the [inet_http_server] section below.
;
; REQUIRES: pip install supervisor
; =============================================================================

; -----------------------------------------------------------------------------
; Unix socket — used by supervisorctl on the local machine
; -----------------------------------------------------------------------------
[unix_http_server]
; Socket lives inside the repo's logs/ directory so it's self-contained.
; %(here)s expands to the directory containing this config file (repo root).
file=%(here)s/logs/supervisor.sock
chmod=0700

; -----------------------------------------------------------------------------
; Web UI — bound to localhost only. Change the password before use.
; -----------------------------------------------------------------------------
[inet_http_server]
; Binds to loopback only — never 0.0.0.0.
port=127.0.0.1:9001

; CHANGE THESE CREDENTIALS before running in any shared environment.
; Username/password are separated by a colon.
; To generate a SHA-1 password hash:
;   python3 -c "import hashlib; print('{SHA}' + __import__('base64').b64encode(hashlib.sha1(b'yourpassword').digest()).decode())"
username=admin
password=changeme

; -----------------------------------------------------------------------------
; supervisord daemon settings
; -----------------------------------------------------------------------------
[supervisord]
; Log the supervisor daemon itself here.
logfile=%(here)s/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info

; PID file lives in logs/ alongside the socket.
pidfile=%(here)s/logs/supervisord.pid

; Do not daemonize — the start script handles backgrounding explicitly.
nodaemon=false

; Strip ANSI colour codes from child process output so log files are clean.
strip_ansi=true

; -----------------------------------------------------------------------------
; supervisorctl — connects via the unix socket above
; -----------------------------------------------------------------------------
[supervisorctl]
serverurl=unix://%(here)s/logs/supervisor.sock

; -----------------------------------------------------------------------------
; RPC interface — required for supervisorctl to work
; -----------------------------------------------------------------------------
[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; =============================================================================
; MANAGED PROCESSES
; =============================================================================

; -----------------------------------------------------------------------------
; ralph — the ceremonies delivery machine
; -----------------------------------------------------------------------------
;
; NOTE ON FOREGROUND EXECUTION:
; `loop.sh start` backgrounds itself and exits immediately, which would cause
; supervisord to restart it in a tight loop. Instead we run the ceremony loop
; inline so supervisord owns the process directly and can manage its lifecycle.
;
; To change the tool (claude → amp), set the RALPH_TOOL environment variable
; in the [program:ralph] environment= line below, or export it before starting
; supervisord.
;
[program:ralph]
; Run the ceremony loop in the foreground. ceremonies.sh exits 0 on success
; (run again) and non-zero on fatal stop (machine needs human attention).
; The outer while-loop mirrors loop.sh's logic but stays in-process.
command=bash -c 'while scripts/ralph/ceremonies.sh --tool "$RALPH_TOOL"; do sleep 5; done'

; Working directory is the repo root — all relative paths in ceremonies.sh
; resolve correctly from here.
directory=%(here)s

; Tool to use. Override by setting RALPH_TOOL in your shell before starting,
; or edit this line. Valid values: claude, amp
environment=RALPH_TOOL="%(ENV_RALPH_TOOL)s"

; Process identity — run as the user who owns the repo, not root.
; Uncomment and set if supervisord itself runs as root (unusual).
; user=youruser

; Restart policy:
; - The while-loop runs ceremonies.sh repeatedly on exit 0 (normal cycle).
; - ceremonies.sh exits 0 even on step failure — it just stops the current
;   cycle early so the while-loop retries from orient on the next iteration.
; - Non-zero exit = actual crash (Python exception, signal). Supervisord
;   restarts the process since autorestart=unexpected.
; - The while-loop breaks on non-zero, so the bash wrapper exits, supervisord
;   sees an unexpected exit, and restarts the wrapper (which re-enters the loop).
autorestart=unexpected
exitcodes=0

; How long the process must run before it's considered "started".
; Ceremonies can take a while to initialise — 10s is conservative.
startsecs=10

; How long to wait for graceful shutdown before force-killing.
; Give ceremonies time to finish the current step cleanly.
stopwaitsecs=30

; SIGINT mirrors Ctrl-C — ceremonies.sh and loop.sh both handle this cleanly.
stopsignal=INT

; Send the stop signal to the entire process group, not just the parent bash
; wrapper. Without these, ceremonies.py (a child of the bash wrapper) never
; receives the signal and supervisord hangs until stopwaitsecs expires.
stopasgroup=true
killasgroup=true

; Redirect stderr into the same log stream as stdout for a single unified log.
redirect_stderr=true

; Stdout log — ceremonies.sh writes its own dated logs under prd/logs/, but
; supervisord captures anything written directly to the process stdout/stderr.
stdout_logfile=%(here)s/logs/ralph.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
